# deploy/systemd/wellwellwell-api.service
[Unit]
Description=WellWellWell API and dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sysadmin
WorkingDirectory=/home/sysadmin/wellwellwell
EnvironmentFile=/home/sysadmin/wellwellwell/.env
ExecStart=/home/sysadmin/wellwellwell/.venv/bin/wellwellwell serve
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
